# backup-system.timer
[Unit]
Description=Backup System daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
