# update-unbound-root-hints.timer
[Unit]
Description=Update Unbound root hints monthly

[Timer]
OnCalendar=monthly
Persistent=true     

[Install]
WantedBy=timers.target
